# templates/default/neutron/neutron-l3-agent.conf
<%= node["admin"]["banner"] %>

description "Quantum l3 plugin agent"
author "Chuck Short <zulcss@ubuntu.com>"

start on runlevel [2345]
stop on runlevel [016]

chdir /var/run

pre-start script
        mkdir -p /var/run/neutron
        chown neutron:root /var/run/neutron
end script

exec start-stop-daemon --start --chuid neutron --exec /usr/bin/neutron-l3-agent -- --config-file=/etc/neutron/neutron.conf --config-file=/etc/neutron/l3_agent_<%= @agent %>.ini --log-file=/var/log/neutron/l3-agent-<%= @agent %>.log
